# search/nmz.cron
#
# Crontab entry for nmz
#
# $Id$
#

SHELL=/bin/sh
MAILTO=root

# Rebuild the indices every hour
17 * * * * $HOME/nmz/search/build-indices.sh
